# deploy/aliyun/his-agent-backend.service
[Unit]
Description=HIS-Agent FastAPI backend
After=network-online.target his-agent-llm-proxy.service
Wants=network-online.target
Requires=his-agent-llm-proxy.service

[Service]
Type=simple
User=hisagent
Group=hisagent
WorkingDirectory=/opt/his-agent/backend
EnvironmentFile=/opt/his-agent/backend/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/his-agent/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1 --log-level info
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
MemoryMax=600M
TasksMax=160

[Install]
WantedBy=multi-user.target
